# api.ini
[program:quarter]
directory=/home/mooc/app/quarter
command=/opt/python3.9.7/bin/gunicorn -c config.py cmd:app
user=mooc
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/quarter.log