# config/api.ini
[program:api]
directory=/var/app/api
command=/opt/python3.10.2/bin/gunicorn -c config.py cmd:app
user=mooc
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/api.log